# backend.conf
[program:taiga-backend]
priority=2
user=taiga
command=start-backend.sh
autostart=true
autorestart=true
directory=/home/taiga/taiga-backend/
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
