# nginx.conf
[program:nginx]
priority=1
command=usr/sbin/service nginx start
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
