# postgresql.conf
[program:postgresql]
priority=1
command=/usr/sbin/service postgresql start
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
